# deploy/etc/s6-overlay/s6-rc.d/server/run
#!/command/with-contenv execlineb

cd /app
foreground { python manage.py createsuperuser --noinput }
python pyapp_docker_example/__main__.py
